# cookbooks/busineme/templates/default/busineme.conf.erb
[program:busineme]
command = <%=@REPODIR%>/gunicorn_script